# deploy/systemd/wizard-defense-support-api.service.example
[Unit]
Description=Wizard Defense Support Preview FastAPI
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/wizard-defense-ai-support
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/ubuntu/wizard-defense-ai-support/.venv/bin/python -m uvicorn backend.app.api:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
